# sources/meta-thermo-mx6/recipes-qt/tfsapp/files/check.service
[Unit]
Description = checking all services 
After = ThermoApp.service

[Service]
ExecStart = /home/root/Scripts/checkservice.sh

[Install]
WantedBy = multi-user.target
